# deploy/dragons-bot.service
[Unit]
Description=Dragons VK Bot (longpoll)
After=dragons-api.service
Wants=dragons-api.service

[Service]
WorkingDirectory=/opt/dragons
Environment=PATH=/opt/dragons/api/venv/bin
Environment=PYTHONPATH=/opt/dragons:/opt/dragons/api
ExecStart=/opt/dragons/api/venv/bin/python -m bot.main
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
